# .config/systemd/user/check-battery.service
[Service]
ExecStart = /home/martin/bin/battery_warning
